# crontab
# Crontab for automatic deployment

# Deploy every hour
0 * * * * $CODE_HOME/deploy.sh
