# crontab
# Crontab for automatic deployment

# Deploy every hour
0 * * * * unbuffer $CODE_HOME/deploy.sh  | ts >> /var/log/cron.txt
